# project/backend/deploy/forus-backend.service
[Unit]
Description=ForUs Backend API
After=network.target

[Service]
Type=simple
User=aloo
WorkingDirectory=/home/aloo/forus/project/backend
EnvironmentFile=/home/aloo/forus/project/backend/.env
ExecStart=/usr/bin/node /home/aloo/forus/project/backend/server.js
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/aloo/forus/project/backend

[Install]
WantedBy=multi-user.target
